# supervisord.conf
[supervisord]
nodaemon=true

[program:nodejs]
directory = /usr/src/iodocs/
command = /bin/bash -c 'export REDIS_URL=$REDIS_PORT && node app.js'
autostart = true
autorestart = true
